# agentic-harness/deploy/agent-studio.service
# Agent Studio orchestrator — /etc/systemd/system/agent-studio.service
# See deploy/vps.md for the full setup and the safety posture this assumes.

[Unit]
Description=Agent Studio orchestrator (studio run --watch)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=studio
Group=studio
WorkingDirectory=/opt/studio/agentic-harness
EnvironmentFile=/opt/studio/agent-studio.env
ExecStart=/opt/studio/agentic-harness/.venv/bin/python -m studio run --watch
Restart=on-failure
RestartSec=30

# Containment: this unit owns its two directories and nothing else.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/studio
PrivateTmp=true

[Install]
WantedBy=multi-user.target
